# transmission-gluetun-port-update/root/etc/s6-overlay/s6-rc.d/svc-mod-transmission-gluetun-port-update/run
#!/usr/bin/with-contenv bash

TRANSMISSION_USERNAME=$USER TRANSMISSION_PASSWORD=$PASS /transmission-gluetun-port-update